# home-manager/files/.config/autostart/xset.desktop
[Desktop Entry]
Name=Xset Commands
Type=Application
Exec=/bin/bash -c "sleep 15 && xset r rate 200 30"
